# tic/lab/QBox/QBox/var/lib/lrpkg/qbox.list
/etc/qbox/*
/etc/init.d/qbox
/sbin/qbox.shape.sh
var/lib/lrpkg/qbox.*

