# tic/lab/Firewall/proxytp/etc/cron.d/mrtg
# /etc/cron.d/mrtg
# $Id: mrtg,v 1.1 2003/10/25 21:42:06 doros Exp $
0-55/5 *	* * *	root	if [ -x /usr/bin/mrtg ] && [ -r /etc/mrtg.cfg ]; then /usr/bin/mrtg /etc/mrtg.cfg >> /var/log/mrtg/mrtg.log 2>&1; fi
